# server/update.crontab
58 23 * * * cd /home/netintech/synthesis/server && bash update.sh backup
58 23 1 * * cd /home/netintech/synthesis/server && bash update.sh clean
*/5 * * * * cd /home/netintech/synthesis/server && bash update.sh check
